# templates/default/graylog2.init.erb
#! /bin/sh

# Installed by Chef => graylog2::default.rb

### BEGIN INIT INFO
# Provides:          graylog2
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the graylog2 server
# Description:       starts the graylog2 java server using the graylog2 control script
### END INIT INFO
#
#  Once installed, use "sudo update-rc.d graylog2 defaults" to enable!

BINDIR=<%= node["graylog2"]["basedir"] %>/server/
NAME=graylog2
DESC=graylog2

test -x /usr/bin/java || exit 0

set -e

start(){
    echo -n "Starting $DESC: "
    cd $BINDIR
    exec sudo java -jar graylog2-server.jar &
    echo "...started."
}
stop() {
    echo -n "Stopping $DESC: "
    sudo ps ax|grep graylog2-server.jar|awk '{print $1}'| xargs sudo /bin/kill  -9  >&2 /dev/null
    echo "...stopped."
}
restart(){
    echo -n "Restarting $DESC: "
    stop
    start
}


case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart}" >&2
        exit 1
        ;;
esac

exit 0
